# service/show_time.timer
[Unit]
Description=Show the time on boot
Requires=show_time.service

[Timer]
Unit=show_time.service
OnBootSec=5

[Install]
WantedBy=timers.target
